# docs/service/exaone-4.0-rag.service
# /etc/systemd/system/exaone-4.0-rag.service
# RAG Mode: ctx=6K, p=1, n=6K
[Unit]
Description=Llama.cpp Server - EXAONE 4.0 32B (RAG Mode, ctx=6K p=1 n=6K)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=kca
WorkingDirectory=/models/llama.cpp
Environment="CUDA_VISIBLE_DEVICES=0"
ExecStart=/models/llama.cpp/build/bin/llama-server \
  -m /models/exaone-4.0/EXAONE-4.0-32B-Q4_K_M.gguf \
  --host 0.0.0.0 \
  --port 8081 \
  --ctx-size 6144 \
  --parallel 1 \
  --cont-batching \
  --defrag-thold 0.1 \
  --jinja \
  -fa auto \
  -ngl 99 \
  --cache-type-k q8_0 \
  --cache-type-v q8_0 \
  --temp 0.6 \
  --top-p 0.95 \
  --repeat-penalty 1.1 \
  --repeat-last-n 256 \
  --presence-penalty 0.3 \
  --frequency-penalty 0.4 \
  --n-predict 6144 \
  --threads 8 \
  --batch-size 4096 \
  --ubatch-size 2048 \
  -a EXAONE-4.0-32B

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=exaone-4.0-rag
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
