# docs/service/llama-server-qwen3vl.service
# /etc/systemd/system/llama-server-qwen3vl.service
[Unit]
Description=Llama.cpp Server - Qwen3-VL 8B Q4_K_M
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=kca
WorkingDirectory=/models/llama.cpp
Environment="CUDA_VISIBLE_DEVICES=0"
ExecStart=/models/llama.cpp/build/bin/llama-server \
  -m /models/qwen3-vl-gguf/Qwen3VL-8B-Instruct-Q4_K_M.gguf \
  --mmproj /models/qwen3-vl-gguf/mmproj-Qwen3VL-8B-Instruct-Q8_0.gguf \
  --host 0.0.0.0 \
  --port 8084 \
  --ctx-size 32768 \
  --parallel 2 \
  --cont-batching \
  --defrag-thold 0.1 \
  --jinja \
  -fa auto \
  -ngl 99 \
  --temp 0.7 \
  --top-p 0.9 \
  --repeat-penalty 1.1 \
  --repeat-last-n 256 \
  --n-predict 2048 \
  --threads 8 \
  --batch-size 2048 \
  --ubatch-size 1024
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llama-server-qwen3vl
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
